# deploy/macos/com.mytho.telegram-claude-bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  launchd LaunchAgent template for the Telegram–Claude bridge.
  Placeholders (__NODE__, __REPO__, __LABEL__) are filled in by install.sh.
  Do not put secrets here — the daemon reads config.json from the repo root.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>__LABEL__</string>

    <key>ProgramArguments</key>
    <array>
        <string>__NODE__</string>
        <string>__REPO__/dist/index.js</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__REPO__</string>

    <!-- Start at login and restart if it ever exits. -->
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <!-- Back off a little on crash loops. -->
    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>__REPO__/data/daemon.out.log</string>
    <key>StandardErrorPath</key>
    <string>__REPO__/data/daemon.err.log</string>

    <!-- launchd starts with a minimal PATH; give node access to common bins. -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
